# root/etc/services.d/sabnzbd/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid sabnzbd /sabnzbd/SABnzbd.py \
	-b 0 -f ${CONFIG} -s ${LISTEN}
